# debian-php7.0/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
echo "[services.d] started cron daemon"
# Runs cron daemon in foreground
cron -l 2 -f
